# etc/profile.d/path.sh
export PATH="${HOME}/.local/bin:${HOME}/.npm/bin:${PATH}"
